# recipes-support/swu-environment/swu-environment/swu-environment.service
[Unit]
Description=SWUpdate Environment

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/bin/swupdate-env"

[Install]
WantedBy=multi-user.target
